# tools/image-builder/files_debootstrap/etc/systemd/system/firecracker.target
[Unit]
Description=Firecracker containerd VM
Requires=basic.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
AllowIsolate=yes
